# n1/files/etc/crontabs/root
0 * * * * "$HOME/scripts/cron_runner.sh" hourly >>"$HOME/logs/hourly.log" 2>&1
0 4 * * * "$HOME/scripts/cron_runner.sh" daily >>"$HOME/logs/daily.log" 2>&1
0 5 * * 1 "$HOME/scripts/cron_runner.sh" weekly >>"$HOME/logs/weekly.log" 2>&1
0 6 1 * * "$HOME/scripts/cron_runner.sh" monthly >>"$HOME/logs/monthly.log" 2>&1
