# n1/files/etc/init.d/cron_reboot
#!/bin/sh /etc/rc.common

START=99

start() {
	if [ -f '/root/scripts/cron_runner.sh' ]; then
		if [ ! -x '/root/scripts/cron_runner.sh' ]; then
			chmod u+x /root/scripts/cron_runner.sh
		fi

		export HOME='/root'
		export USER='root'

		cd /root || return 1

		/root/scripts/cron_runner.sh reboot &
	fi
}
