# n1/files/etc/init.d/cron_reboot
#!/bin/sh /etc/rc.common

START=99

start() {
	task_owner='root'
	runner_script_path="/${task_owner}/scripts/cron_runner.sh"
	export HOME="/${task_owner}"
	export USER="${task_owner}"

	if [ -f "${runner_script_path}" ]; then
		if [ ! -x "${runner_script_path}" ]; then
			chmod u+x "${runner_script_path}"
		fi

		cd "/${task_owner}" || return 1

		"${runner_script_path}" reboot &
	fi

	unset task_owner runner_script_path HOME USER
}
