# n1/files/etc/hotplug.d/iface/99-ipv6-ra
#!/bin/sh

if [ "$ACTION" = 'ifup' ] && [ "$DEVICE" = 'br-lan' ]; then
  if [ -f '/proc/sys/net/ipv6/conf/br-lan/accept_ra' ]; then
    if [ "$(cat /proc/sys/net/ipv6/conf/br-lan/accept_ra)" -ne 2 ]; then
      echo 2 >/proc/sys/net/ipv6/conf/br-lan/accept_ra
    fi
  fi
fi
