# n1/files/etc/hotplug.d/iface/99-ipv6-ra
#!/bin/sh

set -euo pipefail

RA_PATH=/proc/sys/net/ipv6/conf/br-lan/accept_ra

if [ "${ACTION:-}" = ifup ] && [ "${DEVICE:-}" = br-lan ]; then
  [ -f "$RA_PATH" ] || exit 0
  [ "$(cat "$RA_PATH")" -ne 2 ] && echo 2 >"$RA_PATH"
fi
