# n1/files/etc/hotplug.d/iface/99-ipv6-ra
#!/bin/sh

RA_PATH=/proc/sys/net/ipv6/conf/br-lan/accept_ra

main() {
  local ipv6_ra

  [ "${ACTION:-}" = ifup ] && [ "${DEVICE:-}" = br-lan ] || return 0
  [ -f "$RA_PATH" ] && read -rn 1 ipv6_ra <"$RA_PATH" || return 1
  [ "$ipv6_ra" -ne 2 ] && printf '%d\n' 2 >"$RA_PATH"
}

main "$@"
