# n1/files/etc/hotplug.d/iface/99-ipv6-ra
#!/bin/sh

set -euo pipefail

RA_PATH=/proc/sys/net/ipv6/conf/br-lan/accept_ra

main() {
  [ "${ACTION:-}" = ifup ] && [ "${DEVICE:-}" = br-lan ] || return 0
  [ -f "$RA_PATH" ] && [ "$(cat "$RA_PATH")" -ne 2 ] && echo 2 >"$RA_PATH" || :
}

main "$@"
