# n1/files/etc/uci-defaults/99-cron_reboot
#!/bin/sh

if [ -f '/etc/init.d/cron_reboot' ]; then
  if [ ! -x '/etc/init.d/cron_reboot' ]; then
    chmod +x /etc/init.d/cron_reboot
  fi

  /etc/init.d/cron_reboot enable
fi
