# n1/files/etc/uci-defaults/99-cron_reboot
#!/bin/sh

INIT_SERVICE_PATH='/etc/init.d/cron_reboot'

enable_cron_service() {
  service_file="${INIT_SERVICE_PATH}"

  if [ -f "${service_file}" ]; then
    if [ ! -x "${service_file}" ]; then
      chmod +x "${service_file}"
    fi

    "${service_file}" enable
  fi

  unset service_file
}

main() {
  enable_cron_service
}

main "$@"

unset -f enable_cron_service main
unset INIT_SERVICE_PATH
