# n1/files/etc/uci-defaults/99-cron_reboot
#!/bin/sh

set -euo pipefail

[ -f /etc/init.d/cron_reboot ] || exit 0
[ -x /etc/init.d/cron_reboot ] || chmod +x /etc/init.d/cron_reboot

/etc/init.d/cron_reboot enable
