# n1/files/etc/uci-defaults/99-cron_reboot
#!/bin/sh

set -euo pipefail

main() {
  local init_path=/etc/init.d/cron_reboot

  [ -f "$init_path" ] || return 0
  [ -x "$init_path" ] || chmod +x "$init_path"

  "$init_path" enable
}

main "$@"
